# barkdetect.crontab
@reboot /home/pi/barkdetect/startup.sh

# Process recordings every 10 minutes
0-59/10 * * * * bash /home/pi/barkdetect/process.sh 2>> /home/pi/barkdetect/errlog.tx

# Clear junk folder daily
0 6 * * * rm /home/pi/barkdetect/samples/junk/*

# Re-deploy daily
0 6 * * * /home/pi/barkdetect/deploy.sh

# Add a timestamp to the error log daily
0 7 * * * date >> /home/pi/barkdetect/errlog.txt
